# launch_hotspot.timer
[Unit]
Description=Run WiFi check script 20 seconds after boot

[Timer]
OnBootSec=20sec
AccuracySec=1sec
Unit=launch_hotspot.service

[Install]
WantedBy=timers.target
